# cli/telnet@.service
[Unit]
Description=Telnet Server Instance

[Service]
ExecStart=/usr/sbin/in.telnetd
StandardInput=socket